# profile.d/homedir-path.sh
# shellcheck shell=sh

export PATH="$HOME/bin:$HOME/.go/bin:$HOME/.cargo/bin:$PATH:$HOME/.gem/bin:$HOME/.local/bin:$HOME/.composer/vendor/bin"